# docs/example/proxmox/db/mongo-backup.service
# /etc/systemd/system/mongo-backup.service
# Oneshot unit run by mongo-backup.timer. See backup-mongo.sh for install steps.
[Unit]
Description=Logical backup (mongodump) of the open_poll database
After=mongod.service
Wants=mongod.service

[Service]
Type=oneshot
# Credentials (APP_DB_USER / APP_DB_NAME / APP_DB_PASS) are written here by
# db/deploy-db-to-ct.sh (root 600). The leading '-' makes it optional, so a
# manual install that instead edits backup-mongo.sh directly still works.
EnvironmentFile=-/etc/openpoll/mongo-backup.env
ExecStart=/usr/local/bin/backup-mongo.sh
# Don't hold the boot/critical chain if a backup is slow.
Nice=10
IOSchedulingClass=idle
